# BEAM-GUI-DOCKER/beam-gui/app-supervisord.conf
[program:beam_app]
command=/usr/local/bin/launch_app.sh
autostart=true
autorestart=true
startretries=3
priority=50
stdout_logfile=/var/log/app_stdout.log
stderr_logfile=/var/log/app_stderr.log
